# 7_machine_learning/files/jupyterhub.service
[Unit]
Description=Jupyter Hub

[Service]
ExecStart=/usr/bin/jupyterhub
WorkingDirectory=/etc/jupyterhub

[Install]
WantedBy=multi-user.target
